# rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.


cd /www/humidity/ && ./upload.py # We put it in /www/humidity
wifi-live-or-reset

exit 0
